# profile.d/000-kubectl_aliases.sh
source ~/dotfiles/profile.d/.kubectl_aliases
