# profile.d/001-emacs.sh
export EDITOR=emacsclient
alias EMACS='$(which emacs)'
alias es='EMACS -daemon'
alias emacs='emacsclient -n'
